# systemd/user/dock-resume.service
[Unit]
Description=Dock Display Resume Handler (User Session)
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/home/aneesh/.local/bin/dock-resume-handler.sh
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1

[Install]
WantedBy=default.target
